# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:laravel-web]
directory=/app/laravel
command=/bin/sh -c "exec php artisan serve --host=0.0.0.0 --port=${GUI_PORT:-8088}"
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:laravel-queue-worker]
directory=/app/laravel
command=/bin/sh -c "exec php artisan queue:work --sleep=3 --tries=1 --timeout=${QUEUE_WORKER_TIMEOUT:-900}"
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:laravel-worker-job-recovery]
directory=/app/laravel
command=/bin/sh -c "while true; do php artisan worker-jobs:recover --no-interaction || true; sleep ${ARCHIBOT_WORKER_RECOVERY_INTERVAL_SECONDS:-30}; done"
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:event-queue-workers]
directory=/app
command=/bin/sh -c "if [ -n \"${ABSURD_DATABASE_URL:-}\" ] || [ -n \"${DATABASE_URL:-}\" ]; then exec python -m app.event_worker start-workers --concurrency ${ARCHIBOT_QUEUE_WORKER_CONCURRENCY:-1} --claim-timeout ${ARCHIBOT_QUEUE_WORKER_CLAIM_TIMEOUT:-120}; fi; echo 'Event queue worker disabled because neither ABSURD_DATABASE_URL nor DATABASE_URL is configured.'; exec sleep infinity"
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:event-recovery-bridge]
directory=/app
command=/bin/sh -c "if [ -n \"${ABSURD_DATABASE_URL:-}\" ] || [ -n \"${DATABASE_URL:-}\" ]; then exec python -m app.event_worker recovery-scan --interval-seconds ${EVENT_RECOVERY_INTERVAL_SECONDS:-30}; fi; echo 'Event recovery bridge disabled because neither ABSURD_DATABASE_URL nor DATABASE_URL is configured.'; exec sleep infinity"
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:mcp-server]
directory=/app
command=/bin/sh -c "if [ \"${ENABLE_MCP:-false}\" = \"true\" ]; then export MCP_TRANSPORT=${MCP_TRANSPORT:-sse}; exec python -m app.mcp_server; fi; echo 'MCP server disabled because ENABLE_MCP is not true.'; exec sleep infinity"
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
